# overlays/systemd/pandora-gpu-profile.path
[Unit]
Description=Reagir a mudanças de perfil ACPI (desabilitado por padrão — causa loop)
After=graphical-session.target
# NÃO habilitar: PathChanged em platform_profile + escrita no mesmo arquivo
# pelo gpu-profile.sh dispara o service em loop. Use pandora-gpu-profile.timer.

[Path]
PathChanged=/sys/firmware/acpi/platform_profile
PathExists=/sys/firmware/acpi/platform_profile
Unit=pandora-gpu-profile.service

[Install]
WantedBy=graphical-session.target
