# overlays/systemd/pandora-gpu-profile.service
[Unit]
Description=Aplicar perfil de GPU Pandora
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/PandoraProject/scripts/gpu-profile.sh
Environment=PANDORA_ROOT=%h/PandoraProject
